# deploy/systemd/sdwan-bootstrap-agent.service
[Unit]
Description=SD-WAN Bootstrap Agent
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/sdwan-bootstrap-agent --config /etc/sdwan/bootstrap-agent.json
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
